# cocore_installer/cocore.service
[Unit]
Description=CoCore Startup Service

[Service]
Type=simple
ExecStart=/root/init.sh

[Install]
WantedBy=multi-user.target
